# gooderp/root/etc/services.d/gooderp/run
#!/usr/bin/with-contenv sh
gosu good /env/bin/python /app/base/odoo-bin --addons-path=/app/gooderp_addons